# _system/systemd/http-service.example.service
[Unit]
Description=Example hardened HTTP service
After=network.target

[Service]
Type=simple
User=appuser
Group=appgroup
WorkingDirectory=/usr/lib/example-app
EnvironmentFile=-/etc/example-app/env
ExecStart=/usr/lib/example-app/bin/example-app --host 127.0.0.1 --port 8000
Restart=on-failure
UMask=0077
NoNewPrivileges=yes
PrivateTmp=yes
PrivateDevices=yes
DevicePolicy=closed
ProtectSystem=strict
ProtectHome=read-only
ProtectKernelModules=yes
RestrictSUIDSGID=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

[Install]
WantedBy=multi-user.target
