# _system/systemd/scheduled-task.example.service
[Unit]
Description=Example scheduled task service

[Service]
Type=oneshot
User=appuser
Group=appgroup
WorkingDirectory=/usr/lib/example-app
EnvironmentFile=-/etc/example-app/env
ExecStart=/usr/lib/example-app/bin/example-task
UMask=0077
NoNewPrivileges=yes
PrivateTmp=yes
PrivateDevices=yes
DevicePolicy=closed
ProtectSystem=strict
ProtectHome=read-only
ProtectKernelModules=yes
RestrictSUIDSGID=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
